# etc/systemd/system/docker.service.d/custom.conf
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -g /mnt/external/docker -H fd:// --containerd=/run/containerd/containerd.sock
